# templates/flaskapp.service.j2
[Unit]
Description=Flask App
After=network.target

[Service]
Type=simple
User={{ app_user }}
Group={{ app_group }}
WorkingDirectory={{ app_dir }}/current
EnvironmentFile={{ env_file_path }}
ExecStart={{ venv_path }}/bin/gunicorn -w 2 -b 0.0.0.0:5000 app:app
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
